# docker/files/crontab
@daily pacman -Syu --noconfirm && pacman -Scc --noconfirm
