# docker/files/crontab
@daily pacman -Syu --noconfirm && (yes | pacman -Scc)
